# tracks/rootfs/etc/supervisor/conf.d/tracks.conf
[program:tracks]
priority=10
directory=/var/www
command=/var/www/start.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT